# control/scripts/snowstar.service.in
[Unit]
Description=snowstar astrophotography service
After=network-online.target autofs.service
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=1
StartLimitInterval=0
User=@ASTROUSER@
Group=@ASTROGROUP@
Environment="LD_LIBRARY_PATH=@ICELIBDIR@"
WorkingDirectory=@prefix@
ExecStart=@prefix@/sbin/snowstar --debug --base=@prefix@/var/images --database=@prefix@/var/db/database.db --config=@prefix@/var/db/database.db --logfile=@prefix@/var/log/snowstar.log --lines=100000 --files=10 --foreground
StandardOutput=syslog
StandardError=syslog
ExecStartPre=/bin/sleep 10

[Install]
WantedBy=multi-user.target

